# docker/crontab
PATH=/bin:/usr/bin:/usr/local/bin:/opt
#* * * * *  date >/proc/1/fd/1 2>/proc/1/fd/2
*/30 * * * *  /opt/run_leaguesync.sh >/proc/1/fd/1 2>/proc/1/fd/2


 