# docker/crontab
* * * * *  date
*/2 * * * *  /opt/run_leaguesync.sh
